# deploy/amw_analytics-alerts.timer
[Unit]
Description=Run AMW Analytics Notifications & Alerts hourly

[Timer]
OnBootSec=10min
OnUnitActiveSec=1h
Persistent=true
Unit=amw_analytics-alerts.service

[Install]
WantedBy=timers.target
