# deploy/fact-refresh-once.service
[Unit]
Description=TRSM Analytics Fact Refresh (Once)
After=network.target

[Service]
Type=oneshot
WorkingDirectory=/opt/app/amw_analytics
ExecStart=/opt/app/.venv/bin/python run.py refresh-fact --once

[Install]
WantedBy=multi-user.target
