# deploy/trsmanalytics-worker.service
[Unit]
Description=Background Data Loader for trsmanalytics
After=network.target postgresql.service

[Service]
User=kush
Group=kush
WorkingDirectory=/opt/trsmanalytics
Environment="PATH=/opt/trsmanalytics/venv/bin"
Environment="DB_PASS=change_this_strong_password"
Environment="PYTHONUNBUFFERED=1"
ExecStart=/opt/trsmanalytics/venv/bin/python worker.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
