# deploy/trsmanalytics.service
[Unit]
Description=Gunicorn instance to serve trsmanalytics
After=network.target

[Service]
User=kush
Group=www-data
WorkingDirectory=/opt/trsmanalytics
Environment="PATH=/opt/trsmanalytics/venv/bin"
# Ensure you set the correct password here or in an env file
Environment="DB_PASS=change_this_strong_password" 
ExecStart=/opt/trsmanalytics/venv/bin/gunicorn --workers 3 --bind unix:trsmanalytics.sock -m 007 wsgi:app

[Install]
WantedBy=multi-user.target
